# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true

[program:node]
command=node /app/sendContact.js
autorestart=true
